# backup/crontab
*/15 * * * * root /usr/lib/raposfly/backup.sh
